# deploy/omniroam-hostpc.service
[Unit]
Description=OmniRoam HostPC (Go API + C Control Core + Vue UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omniroam
Group=omniroam
WorkingDirectory=/var/lib/omniroam

# 作用：Go 负责 HTTP/API/静态资源，C 控制核心由 -control-core 指向，后续可通过进程或本地 socket 接入。
ExecStart=/usr/sbin/hostpc-api \
  -addr 0.0.0.0:8080 \
  -static /usr/share/omniroam/web/dist \
  -settings /var/lib/omniroam/hostpc-settings.json \
  -users /var/lib/omniroam/hostpc-users.cauth \
  -control-core /usr/lib/omniroam/hostpc-c

Restart=on-failure
RestartSec=3
UMask=0027

[Install]
WantedBy=multi-user.target
